# install/etc/s6/services/10-staytus/run
#!/usr/bin/with-contenv bash

source /assets/functions/*
PROCESS_NAME="staytus"
check_initialized
liftoff

print_info "Starting Puma"
cd /app/staytus
exec bundle exec puma -C config/puma.rb
